# bspwm/.config/systemd/user/float-desktop@.service
[Unit]
Description=Float windows on bspwm desktop

[Service]
ExecStart=%h/.local/bin/float-desktop ^%i

[Install]
WantedBy=graphical.target
